# .config/systemd/user/battery-watcher.service
[Unit]
Description=Personal Battery Notification Service
# We want this to run after the graphical session is ready
After=graphical-session.target

[Service]
# Use %h as a shortcut for your home directory
ExecStart=%h/.local/bin/battery-watcher.sh

# Restart the service if the script ever crashes
Restart=always
RestartSec=10

[Install]
# This "enables" the service to start when you log in
WantedBy=default.target
